# src/tmpl-start-tca.conf
start on runlevel [2345]
console output

script
	exec AGENT_LOC/bin/agent.sh start
end script
