# src/tmpl-stop-tca.conf
start on shutdown
console output

script
	exec AGENT_LOC/bin/agent.sh stop
end script
